# platforms/eyesy_cm3/rootfs/usr/lib/systemd/system/cherrypy.service
[Unit] 
Description=Cherrypy Web Server

[Service] 
Type=simple 
WorkingDirectory=/home/music/EYESY_OS/web
ExecStart=/home/music/EYESY_OS/web/run.sh
User=music
Group=music

[Install] 
WantedBy=multi-user.target
