# platforms/eyesy_cm3/rootfs/etc/systemd/system/eyesysetup.service
[Unit]
Description=EYESY Setup Service
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/etc/eyesysetup.sh
RemainAfterExit=yes
User=root

[Install]
WantedBy=sysinit.target

