# platforms/eyesy_cm3/rootfs/etc/systemd/system/eyesysetup.service
[Unit]
Description=Eyesy Setup Service
After=network.target

[Service]
Type=oneshot
ExecStart=/etc/eyesysetup.sh
RemainAfterExit=yes
User=root

[Install]
WantedBy=multi-user.target

