# platforms/eyesy_cm3/rootfs/etc/systemd/system/fbrestore.service
[Unit]
Description=Restore framebuffer after DRM reset
After=modprobe@drm.service
Requires=modprobe@drm.service

[Service]
#ExecStart=/sbin/greensplash
ExecStart=/usr/bin/cp /dev/random /dev/fb0
StandardOutput=tty
StandardError=journal
Restart=always

[Install]
WantedBy=multi-user.target

