# platforms/eyesy_cm3/rootfs/etc/systemd/system/splashscreen.service
[Unit] 
Description=Splash screen
DefaultDependencies=no
After=local-fs.target

[Service]
ExecStart=/usr/bin/setterm --background cyan --foreground white --clear > /dev/tty1
StandardInput=tty
StandardOutput=tty

[Install]
WantedBy=sysinit.target
