# platforms/eyesy_cm3/rootfs/etc/systemd/system/ttymidi.service
[Unit]
Description=Start ttymidi
After=eyesysetup.service

[Service]
Type=simple
WorkingDirectory=/home/music/
# Run as the 'music' user and group if desired
User=music
Group=music

# Directly run python with unbuffered output
ExecStart=/usr/local/bin/ttymidi -s /dev/ttyAMA0 -b 38400 

# Restart on failure if needed
#Restart=on-failure
#RestartSec=5

[Install]
WantedBy=multi-user.target

