# platforms/eyesy_cm3/rootfs/lib/systemd/system/eyesy-pd.service
[Unit] 
Description=Pd OSC hardware bridge for Eyesy

[Service] 
Type=simple 
WorkingDirectory=/home/music/EYESY_OS/pd
ExecStart=/home/music/EYESY_OS/pd/run.sh
User=music
Group=music

[Install] 
WantedBy=multi-user.target
